# supervisor/supervisord.conf
[unix_http_server]
file=%(ENV_HOME)s/.supervisor/supervisor.sock   ; the path to the socket file
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
logfile=%(ENV_HOME)s/.supervisor/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=%(ENV_HOME)s/.supervisor/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
;umask=022                   ; process file creation umask; default 022
;user=supervisord            ; setuid to this UNIX account at startup; recommended if root
;identifier=supervisor       ; supervisord identifier, default is 'supervisor'
;directory=/tmp              ; default is not to cd during start
;nocleanup=true              ; don't clean up tempfiles at start; default false
;childlogdir=/tmp            ; 'AUTO' child log dir, default $TEMP
;environment=KEY="value"     ; key value pairs to add to environment
;strip_ansi=false            ; strip ansi escape codes in logs; def. false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///%(ENV_HOME)s/.supervisor/supervisor.sock   ; the path to the socket file
;serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

[program:celerybeat]
directory=%(ENV_HOME)s/django_barahlochannel
command=%(ENV_HOME)s/django_barahlochannel/.venv/bin/python -m celery -A celery_barahloch/celery_project --app=celery_barahloch.celery_project.celery_app:app beat --loglevel=info
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stdout_logfile=%(ENV_HOME)s/log/celery/beat.log
stderr_logfile=%(ENV_HOME)s/log/celery/beat.log
;user=fut33v
; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999

[program:celery]
directory=%(ENV_HOME)s/django_barahlochannel
; Set full path to celery program if using virtualenv
command=%(ENV_HOME)s/django_barahlochannel/.venv/bin/python -m celery -A celery_barahloch/celery_project --app=celery_barahloch.celery_project.celery_app:app worker --loglevel=info
numprocs=1
autostart=true
autorestart=true
startsecs=10
stdout_logfile=%(ENV_HOME)s/log/celery/worker.log
stderr_logfile=%(ENV_HOME)s/log/celery/worker.log
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600
; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true
; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
;user=nobody


[program:flower]
numprocs=1
stdout_logfile=%(ENV_HOME)s/log/celery/flower.log
stderr_logfile=%(ENV_HOME)s/log/celery/flower.log
directory=%(ENV_HOME)s/barahl0bot
autostart=true
autorestart=true
startsecs=10
command=python3 -m flower --address=localhost --port=5555 --basic_auth=fut33v:brandnew332XXL
; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true


[program:barahlochannel]
command=python3 %(ENV_HOME)s/django_barahlochannel/barahl0bot/get_photos.py barahlochannel.json 
directory=%(ENV_HOME)s/django_barahlochannel/barahl0bot
autorestart=true

[program:barahlochannel_bot]
command=python3 %(ENV_HOME)s/django_barahlochannel/barahl0bot/barahl0bot.py barahlochannel.json 
directory=%(ENV_HOME)s/django_barahlochannel/barahl0bot
autorestart=true

[program:barahlochannel_mtb]
command=python3 %(ENV_HOME)s/django_barahlochannel/barahl0bot/get_photos.py barahlochannel_mtb.json 
directory=%(ENV_HOME)s/django_barahlochannel/barahl0bot
autorestart=true

[program:barahlochannel_mtb_bot]
command=python3 %(ENV_HOME)s/django_barahlochannel/barahl0bot/barahl0bot.py barahlochannel_mtb.json 
directory=%(ENV_HOME)s/django_barahlochannel/barahl0bot
autorestart=true
